# deploy/cloudflared-r71.service.d/bypass-protonvpn.conf
# Drop-in that reroutes cloudflared-r71.service's traffic around ProtonVPN.
# See deploy/cloudflared-r71-bypass.nft for the nftables rules — they mark
# packets originating from this service's cgroup so the pi's existing
# `from all fwmark 0x80000/0xff0000 lookup 200` ip rule steers them out eth0.
#
# This assumes /etc/cloudflared-r71-bypass.nft has been installed (the nft
# file lives in deploy/ in-repo and is installed alongside the drop-in).
[Service]
ExecStartPost=/usr/sbin/nft -f /etc/cloudflared-r71-bypass.nft
ExecStopPost=-/usr/sbin/nft delete table inet orbit_bypass
