# deploy/cloudflared-r71.service.d/bypass-protonvpn.conf
# Drop-in for cloudflared-r71.service that routes its traffic around the
# ProtonVPN default route. The pi's existing `from all fwmark 0x80000 lookup 200`
# ip rule already sends marked packets out via eth0; this adds an iptables
# rule tagged to cloudflared's systemd cgroup so only its traffic is marked.
[Service]
ExecStartPost=-/bin/sh -c '/usr/sbin/iptables -t mangle -C OUTPUT -m cgroup --path system.slice/cloudflared-r71.service -j MARK --set-mark 0x80000/0xff0000 2>/dev/null || /usr/sbin/iptables -t mangle -A OUTPUT -m cgroup --path system.slice/cloudflared-r71.service -j MARK --set-mark 0x80000/0xff0000'
ExecStopPost=-/usr/sbin/iptables -t mangle -D OUTPUT -m cgroup --path system.slice/cloudflared-r71.service -j MARK --set-mark 0x80000/0xff0000
